# config/systemd/user/llama-stats-bridge.service
[Unit]
Description=llama.cpp + CachyLLama live stats bridge
After=network.target
Wants=network.target

[Service]
Type=simple
ExecStart=/usr/bin/python3 %h/.config/llama.cpp/stats-bridge.py
Restart=on-failure
RestartSec=5
OOMScoreAdjust=-400
MemoryMax=256M
Environment=BRIDGE_PORT=55268
Environment=LLAMA_SERVERS=framearch=http://127.0.0.1:8000,cachy=http://127.0.0.1:9092

[Install]
WantedBy=default.target
