# etc/systemd/system/llama.cpp.service.d/sync-models.conf
[Service]
#ExecStartPre=/usr/local/bin/llama-sync-models.sh

# Gemma 4 multimodal models segfault in Vulkan/RADV when llama-server runs as
# root on Strix Halo. Run the system router as the interactive GPU user; port
# 8000 does not need root, and this matches successful manual launches.
User=juan
Group=juan
SupplementaryGroups=render video audio
Environment=HOME=/home/juan
Environment=XDG_CACHE_HOME=/home/juan/.cache
WorkingDirectory=/home/juan

# Increase timeout: 27B model needs time to unload GPU memory gracefully
TimeoutStopSec=120

# Restart even if main process exits unexpectedly (not just on-failure)
Restart=always
