# etc/systemd/system/llama.cpp.service.d/sync-models.conf
[Service]
#ExecStartPre=/usr/local/bin/llama-sync-models.sh

# Increase timeout: 27B model needs time to unload GPU memory gracefully
TimeoutStopSec=120

# Restart even if main process exits unexpectedly (not just on-failure)
Restart=always
